# conf/saf-configuration-ldap/etc/rc.local
#!/bin/sh -e
#
# rc.local
#

sleep 10s && service nslcd restart

exit 0
